# upstart.conf
start on runlevel [2345]
stop on runlevel [06]

respawn
exec /opt/pe-dhcpd/dhcpd.rb run
